# roles/deadeye/files/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[inet_http_server]
port=*:9001


[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
strip_ansi=true
directory=/

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
history_file=~/.sc_history

[program:deadeyed]
command=/usr/local/bin/deadeyed /usr/local/etc/deadeye.conf
autostart=true
directory=/
